# demo/step1/systemd/system/cats_dogs.target
[Unit]
Description=cats_dogs service
Wants=nginx.service
